# alloy/rootfs/etc/s6-overlay/s6-rc.d/init-alloy/run
#!/usr/bin/with-contenv bash
set -e

CONFIG_DIR="/etc/alloy"
CONFIG_FILE="${CONFIG_DIR}/config.alloy"
OPTIONS_FILE="/data/options.json"

mkdir -p "${CONFIG_DIR}" /data/alloy

# Read options
LOKI_URL=$(jq -r '.loki_url' "${OPTIONS_FILE}")
LOG_LEVEL=$(jq -r '.log_level // "info"' "${OPTIONS_FILE}")
ADDITIONAL_CONFIG=$(jq -r '.additional_config // ""' "${OPTIONS_FILE}")

# Detect journal path
JOURNAL_PATH="/var/log/journal"
if [ ! -d "${JOURNAL_PATH}" ] || [ -z "$(ls -A "${JOURNAL_PATH}" 2>/dev/null)" ]; then
    JOURNAL_PATH="/run/log/journal"
fi

echo "-----------------------------------------------------------"
echo " Grafana Alloy for Home Assistant"
echo " Loki URL: ${LOKI_URL}"
echo " Journal path: ${JOURNAL_PATH}"
echo " Log level: ${LOG_LEVEL}"
echo "-----------------------------------------------------------"

# Generate Alloy config
cat > "${CONFIG_FILE}" <<ALLOYCONFIG
// Grafana Alloy config - auto-generated by HAOS add-on
// Do not edit directly - configure via add-on options

logging {
  level = "${LOG_LEVEL}"
}

// --- Read systemd journal ---
loki.source.journal "journal" {
  path         = "${JOURNAL_PATH}"
  forward_to   = [loki.process.journal.receiver]
  relabel_rules = loki.relabel.journal.rules
  labels       = {
    job = "systemd-journal",
  }
}

// --- Relabel journal fields ---
loki.relabel "journal" {
  forward_to = []

  rule {
    source_labels = ["__journal__systemd_unit"]
    target_label  = "unit"
  }
  rule {
    source_labels = ["__journal__hostname"]
    target_label  = "hostname"
  }
  rule {
    source_labels = ["__journal_syslog_identifier"]
    target_label  = "syslog_identifier"
  }
  rule {
    source_labels = ["__journal__transport"]
    target_label  = "transport"
  }
  rule {
    source_labels = ["__journal_container_name"]
    target_label  = "container_name"
  }
}

// --- Process journal logs ---
loki.process "journal" {
  stage.drop {
    expression = "^\\\\s*$"
  }

  // Extract log level from Home Assistant log format:
  // "2026-03-18 19:02:30.204 INFO (MainThread) [component] message"
  stage.regex {
    expression = "[0-9]{2}:[0-9]{2}:[0-9]{2}[.0-9]* (?P<extracted_level>DEBUG|INFO|NOTICE|WARNING|WARN|ERROR|CRITICAL|FATAL) "
  }

  stage.labels {
    values = {
      level = "extracted_level",
    }
  }

  forward_to = [loki.write.loki.receiver]
}

// --- Write to Loki ---
loki.write "loki" {
  endpoint {
    url = "${LOKI_URL}"
  }
}
ALLOYCONFIG

# Append additional config if provided
if [ -n "${ADDITIONAL_CONFIG}" ] && [ "${ADDITIONAL_CONFIG}" != "null" ]; then
    echo "" >> "${CONFIG_FILE}"
    echo "// --- Additional user config ---" >> "${CONFIG_FILE}"
    echo "${ADDITIONAL_CONFIG}" >> "${CONFIG_FILE}"
fi

echo "Alloy config generated at ${CONFIG_FILE}"
