# alloy/rootfs/etc/s6-overlay/s6-rc.d/init-alloy/run
#!/usr/bin/with-contenv bash
set -e

CONFIG_DIR="/etc/alloy"
CONFIG_FILE="${CONFIG_DIR}/config.alloy"
OPTIONS_FILE="/data/options.json"
GENERATOR="/usr/share/alloy/generate-config.sh"

mkdir -p "${CONFIG_DIR}" /data/alloy

# Read options
export LOKI_URL=$(jq -r '.loki_url // ""' "${OPTIONS_FILE}")
export PROMETHEUS_URL=$(jq -r '.prometheus_url // ""' "${OPTIONS_FILE}")
export PROMETHEUS_USERNAME=$(jq -r '.prometheus_username // ""' "${OPTIONS_FILE}")
export INSTANCE_NAME=$(jq -r '.instance_name // "homeassistant"' "${OPTIONS_FILE}")
export METRICS_SCRAPE_INTERVAL=$(jq -r '.metrics_scrape_interval // "60s"' "${OPTIONS_FILE}")
export LOG_LEVEL=$(jq -r '.log_level // "info"' "${OPTIONS_FILE}")
export ADDITIONAL_CONFIG=$(jq -r '.additional_config // ""' "${OPTIONS_FILE}")

# Require at least one destination
if [ -z "${LOKI_URL}" ] && [ -z "${PROMETHEUS_URL}" ]; then
    echo "FATAL: set at least one of loki_url (logs) or prometheus_url (metrics)." >&2
    exit 1
fi

# Detect journal path (logs only)
export JOURNAL_PATH="/var/log/journal"
if [ ! -d "${JOURNAL_PATH}" ] || [ -z "$(ls -A "${JOURNAL_PATH}" 2>/dev/null)" ]; then
    export JOURNAL_PATH="/run/log/journal"
fi

echo "-----------------------------------------------------------"
echo " Grafana Alloy for Home Assistant"
[ -n "${LOKI_URL}" ]       && echo " Logs    -> Loki:       ${LOKI_URL}"
[ -n "${PROMETHEUS_URL}" ] && echo " Metrics -> Prometheus: ${PROMETHEUS_URL} (every ${METRICS_SCRAPE_INTERVAL})"
[ -n "${LOKI_URL}" ]       && echo " Journal path: ${JOURNAL_PATH}"
echo " Instance: ${INSTANCE_NAME}"
echo " Log level: ${LOG_LEVEL}"
echo "-----------------------------------------------------------"

# Generate Alloy config
"${GENERATOR}" > "${CONFIG_FILE}"

echo "Alloy config generated at ${CONFIG_FILE}"
